# src/app/service/sigaa/run
rm -r data 
rm -r debug

mkdir data
mkdir debug

../../../../venv/bin/python3 get_components_list.py
