# src/app/service/sigaa/run
rm data/* && rm debug/* && ../../../../venv/bin/python3 get_components_list.py
